# remote/notify.cron
00 07 * * * root /root/bday/remote/notify.sh